# automated_install/configs/supervisor/laravel_echo_server.conf
[program:laravel-echo-server]
process_name=%(program_name)s_%(process_num)02d
directory=/var/www/html/giantbomb_downloader/
command=laravel-echo-server start
autostart=true
autorestart=true
numprocs=1
#redirect_stderr=true
#stdout_logfile=/var/www/html/giantbomb_downloader/echo-server-worker.log
