# supervisor/laravel_echo_server.conf
[program:laravel-echo-server]
process_name=%(program_name)s_%(process_num)02d
command=php /home/forge/app.com/laravel-echo-server start
autostart=true
autorestart=true
user=forge
numprocs=1
redirect_stderr=true
stdout_logfile=/home/forge/app.com/worker.log
